# ceilo-agent-notification.conf
[program:ceilometer-agent-notification]
command=/usr/bin/start-ceilo-agent-notification.sh
autostart=true
autorestart=false
